# ops/roscoe-slack.service
[Unit]
Description=Roscoe Slack Socket-Mode Bridge
After=network-online.target
Wants=network-online.target

[Service]
# Update WorkingDirectory/EnvironmentFile to match your VM path.
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/roscoe
EnvironmentFile=/opt/roscoe/.env
ExecStart=/usr/bin/env bash scripts/run_slack_bot.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target

